# trading-data-collector.service
[Unit]
Description=Trading Data Collector
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-api
Environment=PATH=/home/ubuntu/trading-api/vps_env/bin
ExecStart=/home/ubuntu/trading-api/vps_env/bin/python data_collector.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
