# trading-api-server.service
[Unit]
Description=Trading API Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-api
Environment=PATH=/home/ubuntu/trading-api/vps_env/bin
ExecStart=/home/ubuntu/trading-api/vps_env/bin/gunicorn --bind 0.0.0.0:5000 --workers 4 --timeout 120 api_server:app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
